# elasticsearch/supervisord.conf
[supervisord]
nodaemon=true
logfile = /tmp/supervisord.log
loglevel = DEBUG

[program:elasticsearch]
command=/elasticsearch/bin/elasticsearch
user=elasticsearch
