# cassandra/supervisord.conf
[supervisord]
nodaemon=true
logfile = /tmp/supervisord.log
loglevel = DEBUG

[program:cassandra]
command=/etc/cassandra/bin/cassandra -f
user=cassandra
